# aliases/profile.d/aliases_basic
# .aliases file for shell aliases & related envars

export GHE_HOME=~/ghe-pip2
alias cd-ghe='cd $GHE_HOME'

export GH_HOME=~/github
alias cd-gh='cd $GH_HOME'
alias cd-dev-env='cd-gh && cd dev-env'

alias l='ls -la'

