# fzf/profile.d/fzf
# Adds necessary fzf sourcing line
[ -f ~/.fzf.bash ] && source ~/.fzf.bash
[ -f ~/.fzf.zsh ] && source ~/.fzf.zsh

export FZF_DEFAULT_COMMAND='fd . --hidden --exclude ".git/*"'
