# ssh/profile.d/ssh_tunnel
# Function to establish tunnel to nodes in Controls Network
# NOTE: Requires the ssh setup in [this gist](https://gist.ghe-pip2.fnal.gov/dcrisp/73f913b2907f1a81e86a3bbae2cfe240)
function tunnel() {
	local node="$([[ -p /dev/stdin ]] && cat - || echo "$@")"
	[[ -z "$node" ]] && return 1
	echo Connecting to $node
	ssh -L 6000:localhost:6000 $node
}
