# nvim/profile.d/bob_nvim
BOB_NVIM_PATH="$HOME/.local/share/bob/nvim-bin"

export PATH=$BOB_NVIM_PATH:$PATH
