# shell-aliases/profile.d/aliases_basic
# .aliases file for shell aliases & related envars

export GHE_HOME="~/ghe-pip2"
alias cdghe='cd $GHE_HOME'

export GH_HOME="~/github"
alias cdgh='cd $GH_HOME'

alias l='ls -a'
alias ll='ls -la'
